# roles/catalogue/files/systemd.service
[Unit]
Description=Catalogue Service

[Service]
ExecStart=/usr/bin/node /home/catalogue/server.js

[Install]
WantedBy=multi-user.target